# counter/usr/local/bin/counter.service
[Unit]
Description=Simple Counter Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/counter_service.py
Restart=always
User=counter
Group=counter
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target



